# services/indexation/crontab
57 17 * * * /usr/local/bin/python /index.py >> cron_log 2>&1
